# trunk/Kettle8/CALL8/scripts/main/bin/.cron
# Format:
# minutes_of_hour hours_of_day days_of_month months_of_year days_of_the_week command_to_execute
#Schedule to run once on 16th of every month at 6.00am
0 6 * * * /u01/maximus/maxdat-uat/CALL8/scripts/main/bin/run_call_provider_load.sh > /u01/maximus/maxdat-uat/CALL8/logs/cron_call_prov_load_$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1

