# ContactCenter/trunk/kettle/MAXDAT/implementation/DCHIX/DCHIX/scripts/ContactCenter/main/bin/.cron
#COEESMAP-1911
#0 4-22 * * * /u01/maximus/maxdat-prd/CO/ETL/scripts/coats_run_bpm.sh 2>&1|sed "s/^/$(date)/" >> /u01/maximus/maxdat-prd/CO/ETL/logs/coats_run_bpm.log
#COEESMAP-1911
#0 4-22 * * * /u01/maximus/maxdat-prd/CO/ETL/scripts/coats_run_planning.sh 2>&1|sed "s/^/$(date)/" >> /u01/maximus/maxdat-prd/CO/ETL/logs/coats_run_planning.log
#COEESMAP-1937
#10,25,40,55 * * * * /u01/maximus/maxdat-prd/CO/ETL/scripts/coats_run_planning.sh 2>&1|sed "s/^/$(date)/" >> /u01/maximus/maxdat-prd/CO/ETL/logs/coats_run_planning.log
#15 * * * * /u01/maximus/maxdat-prd/CO/ETL/scripts/ContactCenter/main/bin/manage_scheduled_contact_center_jobs.sh > /u01/maximus/maxdat-prd/CO/ETL/logs/ContactCenter/`uname -n`_contcent.log 2>&1
#20  * * * * /u01/maximus/maxdat-prd/CO/ETL/scripts/ContactCenter/main/bin/manage_adhoc_contact_center_jobs.sh > /u01/maximus/maxdat-prd/CO/ETL/logs/ContactCenter/`uname -n`_contcent.log 2>&1 
#00 06 * * * /u01/maximus/maxdat-prd/CO/ETL/scripts/ContactCenter/main/bin/load_New_Queue_Flag_for_Alert.sh > /u01/maximus/maxdat-prd/CO/ETL/logs/ContactCenter/`uname -n`_Alert_Flag 2>&1
#MAXDAT-2127
15 * * * * sh /u01/maximus/maxdat-prd/MDHIX/ETL/scripts/ContactCenter/main/bin/manage_adhoc_contact_center_jobs.sh > /u01/maximus/maxdat-prd/MDHIX/ETL/logs/ContactCenter/`uname -n`_contcent.log 2>&1
0 * * * * sh /u01/maximus/maxdat-prd/MDHIX/ETL/scripts/ContactCenter/main/bin/manage_scheduled_contact_center_jobs.sh > /u01/maximus/maxdat-prd/MDHIX/ETL/logs/ContactCenter/`uname -n`_contcent.log 2>&1
#MAXDAT-2413 at 6 pm  daily
0 18 * * * /u01/maximus/maxdat-prd/MDHIX/ETL/scripts/Run_PP_WFM_RUNALL.sh >> /u01/maximus/maxdat-prd/MDHIX/ETL/logs/Run_PP_WFM_RUNALL.log &
15 * * * * /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/scripts/ContactCenter/main/bin/manage_adhoc_contact_center_jobs.sh >> /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/manage_adhoc_contact_center_jobs.log 2>&1
#MAXDAT-3760
2 * * * * /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/scripts/ContactCenter/main/bin/scheduled_contact_center_job_executions.sh >> /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/scheduled_contact_center_job_executions.log
#MAXDAT-3498
15 0-23 * * * /u01/maximus/maxdat-prd/KSCH_Legacy/scripts/ContactCenter/main/bin/manage_adhoc_contact_center_jobs.sh >> /u01/maximus/maxdat-prd/KSCH_Legacy/logs/manage_adhoc_contact_center_jobs.log &
#MAXDAT-3611
0 * * * * /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/scripts/ContactCenter/implementation/CiscoEnterprise/bin/Run_Daily_IVR_Files.sh > /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/ContactCenter/Run_Daily_IVR_Files_$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1
#MAXDAT-3720
0 7 * * 1 /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/scripts/ContactCenter/main/bin/run_scheduled_export_amp_metrics.sh "WEEKLY" > /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/ContactCenter/cron_run_weekly_export_amp_metrics_$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1
5 7 2 * * /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/scripts/ContactCenter/main/bin/run_scheduled_export_amp_metrics.sh "MONTHLY" > /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/ContactCenter/cron_run_monthly_export_amp_metrics_$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>&1
#MAXDAT-3882
30 17 2 8 2 /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/scripts/ContactCenter/main/bin/run_scheduled_export_amp_metrics.sh "WEEKLY" > /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/ContactCenter/cron_run_weekly_export_amp_metrics_$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>/u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/ContactCenter/cron_error_weekly.log 
30 16 2 8 2 /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/scripts/ContactCenter/main/bin/run_scheduled_export_amp_metrics.sh "MONTHLY" > /u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/ContactCenter/cron_run_monthly_export_amp_metrics_$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>/u01/maximus/maxdat-prd/CiscoEnterprise/ETL/logs/ContactCenter/cron_error_monthly.log 
#MAXDAT-3906
0 7 16 * * /u01/maximus/maxdat-prd/AO/ETL/main/bin/run_ao_input_jobs.sh > /u01/maximus/maxdat-prd/AO/logs/cron_ao_input_$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>/u01/maximus/maxdat-prd/AO/logs/cron_error.log 
*/30 * 17,18,19,20,21,22,23,24,25 * * /u01/maximus/maxdat-prd/AO/ETL/main/bin/run_ao_extract_jobs.sh > /u01/maximus/maxdat-prd/AO/logs/cron_ao_input_$(date +\%Y-\%m-\%d_\%H-\%M-\%S).log 2>/u01/maximus/maxdat-prd/AO/logs/cron_error.log
