# ContactCenter/trunk/kettle/MAXDAT/implementation/LAEB/LAEB/scripts/ContactCenter/main/bin/.cron
# Format:
# minutes_of_hour hours_of_day days_of_month months_of_year days_of_the_week command_to_execute
00 05,06,07,08,09,10,11,12,13,14,15,16,17,18,19,20 * * * /u01/maximus/maxdat-uat/TNRD/scripts/cronfiles/cron_tn_run_bpm.sh 2>&1
#manage_adhoc_contact_center_jobs.sh cron job to run every hour every day on the 15th minute of the hour 
15 * * * * /u01/maximus/maxdat-uat/FOLSOM/ETL/scripts/ContactCenter/main/bin/manage_adhoc_contact_center_jobs.sh > /u01/maximus/maxdat-uat/FOLSOM/ETL/logs/manage_adhoc_contact_center_jobs.log 2>&1
#cron job to run every hour every day on the top of the hour 
00 * * * * /u01/maximus/maxdat-uat/FOLSOM/ETL/scripts/ContactCenter/main/bin/manage_scheduled_contact_center_jobs.sh > /u01/maximus/maxdat-uat/FOLSOM/ETL/logs/manage_scheduled_contact_center_jobs.log 2>&1

